# supervisord.conf
[program:node]
command = node server.js
directory = /home/dotcloud/current

[program:dpwolfnodechat]
command = node /home/dotcloud/current/server.js